# Resources/bin/z-proxy
#!/bin/sh /etc/rc.common
START=96
STOP=10

PROG=/usr/local/bin/proxy-z



reload() {
	stop
    sleep 1
    start
	return 0
}

restart() {
	stop
    sleep 1
    start
}

start() {
    sysctl -w net.core.rmem_max=21299200;
    # HOME="/root" /usr/local/bin/proxy-z -http -d  2>&1  1>/var/log/proxy-z.log;
    /usr/sbin/auto-z-proxy
    
}

stop() {
	ps | grep $PROG | grep -v grep | awk '{print $1}' | xargs kill -9
    ps | grep auto-z-proxy | grep -v grep | awk '{print $1}' | xargs kill -9
	return 0
}